# scripts/monitoring/sutazaiapp-monitor.service
[Unit]
Description=Sutazaiapp Resource Monitor
After=docker.service
Wants=docker.service

[Service]
Type=simple
Restart=always
RestartSec=10
User=root
WorkingDirectory=/opt/sutazaiapp
Environment="PYTHONPATH=/opt/sutazaiapp"
Environment="MEMORY_THRESHOLD=75"
Environment="CPU_THRESHOLD=70"
Environment="CHECK_INTERVAL=60"
ExecStart=/usr/bin/python3 /opt/sutazaiapp/scripts/monitoring/automated_resource_monitor.py --daemon
StandardOutput=append:/var/log/sutazaiapp-monitor.log
StandardError=append:/var/log/sutazaiapp-monitor-error.log

# Resource limits for the monitor itself
MemoryLimit=128M
CPUQuota=10%

# Restart on failure
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target