# scripts/systemd/sutazai-maintenance.timer
[Unit]
Description=Run SutazAI System Maintenance Daily

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target 