# ultracleanup_backup_20250811_145137/agents/hardware-resource-optimizer/tests/hardware-optimizer-tests.service
[Unit]
Description=Hardware Optimizer Continuous Test Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/sutazaiapp/agents/hardware-resource-optimizer/tests
ExecStart=/usr/bin/python3 /opt/sutazaiapp/agents/hardware-resource-optimizer/tests/continuous_validator.py --interval 60 --dashboard
Restart=always
RestartSec=10
StandardOutput=append:/var/log/hardware-optimizer-tests.log
StandardError=append:/var/log/hardware-optimizer-tests.log

# Environment
Environment="PYTHONUNBUFFERED=1"

# Resource limits
LimitNOFILE=65536
MemoryLimit=512M
CPUQuota=25%

[Install]
WantedBy=multi-user.target