# backend/systemd/sutazai-memory-optimizer.service
[Unit]
Description=SutazAI Memory Optimizer Service
After=network.target

[Service]
Type=oneshot
User=root
ExecStart=/opt/sutazaiapp/scripts/memory_optimizer.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sutazai-memory-optimizer

[Install]
WantedBy=multi-user.target 