# scripts/mcp/mcp-cleanup.service
[Unit]
Description=MCP Container Cleanup Daemon
Documentation=https://github.com/sutazai/mcp-container-management
After=docker.service
Requires=docker.service
StartLimitBurst=3
StartLimitIntervalSec=60

[Service]
Type=simple
User=root
Group=docker
ExecStart=/opt/sutazaiapp/scripts/mcp/cleanup_containers.sh --daemon
ExecStop=/bin/kill -TERM $MAINPID
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10
TimeoutStartSec=60
TimeoutStopSec=30

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/tmp /opt/sutazaiapp/logs
PrivateTmp=true

# Resource limits
LimitNOFILE=1024
LimitNPROC=100
MemoryMax=128M

# Environment
Environment=MCP_CLEANUP_INTERVAL=300
Environment=MCP_MAX_AGE=3600
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mcp-cleanup

[Install]
WantedBy=multi-user.target