# systemd/sutazai-cpu-monitor.service
[Unit]
Description=SutazAI CPU Usage Monitor Service
After=network.target

[Service]
Type=simple
User=root
ExecStart=/opt/sutazaiapp/scripts/cpu_monitor.sh
Restart=on-failure
RestartSec=30
Nice=19
IOSchedulingClass=idle
CPUSchedulingPolicy=idle
CPUSchedulingPriority=19

[Install]
WantedBy=multi-user.target 