# systemd/sutazai-memory-optimizer.timer
[Unit]
Description=SutazAI Memory Optimizer Timer
After=network.target

[Timer]
# Run every 30 minutes
OnUnitActiveSec=30m
OnBootSec=15m
AccuracySec=1m

[Install]
WantedBy=timers.target 