# systemd/sutazaiapp.service
[Unit]
Description=SutazaiApp Service
After=network.target

[Service]
User=sutazaidev
Group=sutazaidev
WorkingDirectory=/opt/sutazaiapp
Environment="PATH=/opt/venv-sutazaiapp/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONPATH=/opt/sutazaiapp:/opt/venv-sutazaiapp/lib/python3.11/site-packages"
ExecStart=/opt/venv-sutazaiapp/bin/python3.11 -m uvicorn backend.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target