# systemd/sutazai-grafana.service
[Unit]
Description=SutazAI Grafana Dashboards Service
After=network.target sutazai-prometheus.service
Requires=sutazai-prometheus.service

[Service]
Type=simple
User=sutazaiapp
Group=sutazaiapp
WorkingDirectory=/opt/sutazaiapp
ExecStart=/opt/sutazaiapp/scripts/start_grafana.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sutazai-grafana
TimeoutStartSec=90
LimitNOFILE=65536

# Resource controls
CPUQuota=30%
MemoryLimit=1G

[Install]
WantedBy=multi-user.target 