# _legacy/sutazai-real-agents.service
[Unit]
Description=SutazAI Real Agent System
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/sutazaiapp
ExecStart=/opt/sutazaiapp/start_sutazai_with_real_agents.sh
ExecStop=/usr/bin/docker-compose -f /opt/sutazaiapp/docker-compose-real-agents.yml down
User=root
Group=root

[Install]
WantedBy=multi-user.target