# scripts/maintenance/memory-monitor.service
[Unit]
Description=Memory Bank Monitor and Cleanup Service
Documentation=file:///opt/sutazaiapp/docs/memory-cleanup.md
After=network.target

[Service]
Type=oneshot
ExecStart=/opt/sutazaiapp/scripts/maintenance/memory_monitor.sh
StandardOutput=journal
StandardError=journal
User=root
Group=opt-admins

# Security settings
PrivateTmp=yes
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/sutazaiapp/memory-bank /var/log

[Install]
WantedBy=multi-user.target