# scripts/maintenance/memory-monitor.timer
[Unit]
Description=Run Memory Bank Monitor every hour
Documentation=file:///opt/sutazaiapp/docs/memory-cleanup.md
Requires=memory-monitor.service

[Timer]
# Run every hour
OnCalendar=hourly
# Run 5 minutes after boot
OnBootSec=5min
# If missed, run immediately
Persistent=true

[Install]
WantedBy=timers.target