# backend_archive/additional_cleanup/misc_files/systemd/sutazai-orchestrator.service
[Unit]
Description=SutazaiApp AI Orchestrator Service
After=network.target sutazaiapp.service

[Service]
User=sutazaidev
Group=sutazaidev
WorkingDirectory=/opt/sutazaiapp
Environment="PYTHONPATH=/opt/sutazaiapp"
ExecStart=/opt/venv-sutazaiapp/bin/python3.11 -m ai_agents.orchestrator.orchestrator
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target 