# backend_archive/additional_cleanup/misc_files/systemd/sutazai-backend.service
[Unit]
Description=SutazAI Backend Service
After=network.target
StartLimitIntervalSec=60

[Service]
Type=simple
User=sutazaidev
Group=sutazaidev
WorkingDirectory=/opt/sutazaiapp
Environment=PYTHONPATH=/opt/sutazaiapp
ExecStart=/opt/sutazaiapp/venv/bin/uvicorn backend.main:app --host 0.0.0.0 --port 8000 --workers 4
Restart=always
RestartSec=5
StandardOutput=append:/opt/sutazaiapp/logs/backend/backend.log
StandardError=append:/opt/sutazaiapp/logs/backend/backend.error.log

[Install]
WantedBy=multi-user.target 