# backend_archive/additional_cleanup/misc_files/systemd/sutazai.service
[Unit]
Description=SutazAI Application Service
After=network.target postgresql.service
Wants=network.target

[Service]
Type=forking
User=sutazaidev
Group=sutazaidev
WorkingDirectory=/opt/sutazaiapp
ExecStart=/opt/sutazaiapp/scripts/start_all.sh
ExecStop=/opt/sutazaiapp/scripts/stop_all.sh
TimeoutStartSec=300
TimeoutStopSec=120
Restart=on-failure
RestartSec=30
# Setup auto-maintenance service to run 4 hours after startup
ExecStartPost=/bin/bash -c 'sleep 14400 && /opt/sutazaiapp/scripts/auto_maintenance.sh --cron'

[Install]
WantedBy=multi-user.target 