# scripts/supreme-ai-health-monitor.service
[Unit]
Description=Supreme AI Health Monitor Service
After=supreme-ai-orchestrator.service
Wants=supreme-ai-orchestrator.service

[Service]
Type=simple
User=sutazaiapp_dev
Group=sutazaiapp_dev
WorkingDirectory=/opt/sutazaiapp
Environment=PYTHONPATH=/opt/sutazaiapp
ExecStart=/opt/sutazaiapp/venv/bin/python3.11 /opt/sutazaiapp/scripts/monitor_health.py
Restart=always
RestartSec=10
StandardOutput=append:/opt/sutazaiapp/logs/health_monitor.log
StandardError=append:/opt/sutazaiapp/logs/health_monitor.error.log

[Install]
WantedBy=multi-user.target 