# scripts/supreme-ai-orchestrator.service
[Unit]
Description=Supreme AI Orchestrator Service
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=sutazaiapp_dev
Group=sutazaiapp_dev
WorkingDirectory=/opt/sutazaiapp
Environment=PYTHONPATH=/opt/sutazaiapp
ExecStart=/opt/sutazaiapp/scripts/manage_orchestrator.sh start
ExecStop=/opt/sutazaiapp/scripts/manage_orchestrator.sh stop
Restart=always
RestartSec=10
StandardOutput=append:/opt/sutazaiapp/logs/orchestrator.log
StandardError=append:/opt/sutazaiapp/logs/orchestrator.error.log

[Install]
WantedBy=multi-user.target 