# scripts/sutazai-monitor.service
[Unit]
Description=SutazAI Health Monitor Service
After=docker.service
Requires=docker.service

[Service]
Type=simple
Restart=always
RestartSec=30
User=root
WorkingDirectory=/opt/sutazaiapp
ExecStart=/usr/bin/python3 /opt/sutazaiapp/scripts/health_monitor.py --continuous
StandardOutput=append:/opt/sutazaiapp/logs/monitor.log
StandardError=append:/opt/sutazaiapp/logs/monitor.log

[Install]
WantedBy=multi-user.target
